# systemd/inetmonitor.service
[Unit]
Description=Monitor Internet connection
[Service]
ExecStart=/opt/inetmonitor/bin/pingone
[Install]
WantedBy=timers.target


